# scripts/oca_nats_config_update.service
[Unit]
Description=Run auto_update.sh every 5 minutes

[Service]
ExecStart=/usr/bin/auto_update.sh